# scripts/linux/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty 115200 console vt100
::restart:/sbin/init
::ctrlaltdel:/etc/init.d/rc6
